# examples/systemd/portreach-agent.service
[Unit]
Description=portreach agent (network reachability probe)
Documentation=https://github.com/lavr/portreach
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/local/bin/portreach agent --listen :8732
Environment=NODE_NAME=%H
DynamicUser=yes
Restart=on-failure
RestartSec=2
AmbientCapabilities=
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
